# MultiAgent/PressRelease/crontab
0 */2 * * * cd /app && /usr/local/bin/python pipeline.py --source rbi > /proc/1/fd/1 2>&1

